# agent/systemd/vps-insight.timer
[Unit]
Description=Run VPS Insight agent every 10 seconds

[Timer]
OnBootSec=30s
OnUnitActiveSec=10s
AccuracySec=1s
Persistent=true

[Install]
WantedBy=timers.target
